# supervisord.conf
[supervisord]
pidfile=/run/supervisor.pid
logfile=/var/log/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
nodaemon=true                ; (start in foreground if true;default false)
user=root

[inet_http_server]
port=127.0.0.1:9001

[unix_http_server]
file=/run/supervisor.sock   ; the path to the socket file

[supervisorctl]
serverurl=unix:///run/supervisor.sock ; use a unix:// URL  for a unix socket

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:concierge]
directory=/app
command=/usr/local/bin/python main.py
autostart=true
autorestart=true
priority=5
stdout_events_enabled=true
stderr_events_enabled=true
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
user=concierge

[eventlistener:concierge_exit]
command=/app/kill.py
events=PROCESS_STATE_FATAL,PROCESS_STATE_STOPPED

[include]
files = /app/tasks.d/*.conf
